# roles/external_backup/templates/zfs-backup@.service.j2
[Unit]
Description=ZFS Backup Job for Drive ID %i

[Service]   
Type=oneshot
ExecStartPre=/usr/bin/test -e /dev/disk/by-id/ata-%i
ExecStart=/usr/local/bin/zfs-backup.sh %i